# prpl-demo-utils/files/etc/uci-defaults/90-demo-adjustments
uci -q set dsl.line.type=intelcpe
uci -q commit dsl

uci -q set network.lan.is_lan=1
uci -q set network.lan_dev.vid=1
uci -q set network.wan_dev.vid=2
uci -q commit network

uci -q set wireless.status=wifi-status
uci -q set wireless.status.wlan=1
uci -q set wireless.wlan0.band=b
uci -q set wireless.wlan2.band=a
uci -q commit wireless

uci -q set firewall.settings=settings
uci -q set firewall.settings.disabled=0
uci -q set firewall

