# prpl-demo-utils/files/etc/uci-defaults/90-demo-adjustments
uci -q set dsl.line.type=intelcpe
uci -q commit dsl

uci -q set network.wan.ifname="ptm0.10"
uci -q delete network.lan_dev
uci -q set network.wan.macaddr=$(uci -q get network.wan_dev.macaddr)
uci -q delete network.wan_dev
uci -q set network.lan.is_lan=1
uci -q set network.wan6.ifname="@wan"
uci -q commit network

uci -q set wireless.status=wifi-status
uci -q set wireless.status.wlan=1
uci -q set wireless.wlan0.band=b
uci -q set wireless.wlan2.band=a
uci -q commit wireless

uci -q set firewall.settings=settings
uci -q set firewall.settings.disabled=0
uci -q commit firewall

uci -q add firewall rule
uci -q set firewall.@rule[-1].name='allow-video'
uci -q set firewall.@rule[-1].src='wan'
uci -q set firewall.@rule[-1].dest='lan'
uci -q set firewall.@rule[-1].dest_ip='224.0.0.0/4'
uci -q set firewall.@rule[-1].proto='udp'
uci -q set firewall.@rule[-1].family='ipv4'
uci -q set firewall.@rule[-1].target='ACCEPT'
uci -q commit firewall

uci -q set system.@system[0].hostname=prplWrt
uci -q commit system
