# prpl-demo-utils/files/etc/uci-defaults/90-intel-adjustments
uci -q set dsl.line.type=intelcpe
uci -q commit dsl

