# prpl-demo-utils/files/etc/hotplug.d/iface/90-online
[ "$INTERFACE" == "wan" ] || exit

if [ "$ACTION" == "ifup" ]; then
	echo "{ \"online\" : true }" > /tmp/internet_connection_status
	ubus send internet '{"status" : "online"}'
elif [ "$ACTION" == "down" ]; then
	echo "{ \"online\" : false }" > /tmp/internet_connection_status
	ubus send internet '{"status" : "offline"}'
fi

